# files/systemd/backups.timer
[Unit]
Description=NAS backups timer

[Timer]
OnCalendar=Tue,Thu,Fri 10:00:00

[Install]
WantedBy=timers.target
