# roles/backups/templates/backups.service.j2
[Unit]
Description=NAS backups

[Service]
Type=oneshot
ExecStart=/home/{{ ansible_user }}/bin/backups
User={{ ansible_user }}
