# roles/backups/files/backups.service
[Unit]
Description=NAS backups

[Service]
Type=oneshot
ExecStart=/home/jarias/.bin/backups
User=jarias
