# roles/backups/files/tarsnap.service
[Unit]
Description=Tarsnap backups

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tarsnap.sh
