# roles/backups/files/tarsnap.service
[Unit]
Description=Tarsnap backups

[Service]
Type=oneshot
ExecStart=/home/jarias/tarsnap.sh
