# roles/backups/files/tarsnap.service
[Unit]
Description=Tarsnap backups

[Service]
Type=oneshot
ExecStart=/root/tarsnap.sh
